# agent/digital.keiyi.idle_agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>digital.keiyi.idle_agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>AGENT_PATH_PLACEHOLDER/keiyi_agent_app.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>LOG_PATH_PLACEHOLDER/keiyi_idle_agent_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>LOG_PATH_PLACEHOLDER/keiyi_idle_agent_stderr.log</string>

    <key>WorkingDirectory</key>
    <string>AGENT_PATH_PLACEHOLDER</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
